# oracle/blockchain-lab.service
[Unit]
Description=Blockchain Lab - Educational Blockchain Network
Documentation=https://github.com/brighammurdoch-byte/The-Blockchain-Lab
After=network.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/apps/blockchain-demo
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /home/opc/apps/blockchain-demo/bin/www
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blockchain-lab

# Security hardening (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
